# usr/lib/systemd/user/render-menu-options-user.path
[Unit]
Description=Add rendering type in menu

[Path]
PathModified=%h/.local/share/applications
Unit=render-menu-options-user.service

[Install]
WantedBy=default.target
