# usr/lib/systemd/system/render-menu-options.path
[Unit]
Description=Add rendering type in menu

[Path]
PathModified=/usr/share/applications
Unit=render-menu-options.service

[Install]
WantedBy=default.target
